# .docker/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=20

[program:queue-worker]
command=php artisan queue:work --tries=3 --timeout=60
directory=/var/www
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/logs/queue.log
stderr_logfile=/var/www/storage/logs/queue-error.log
priority=30

[program:scheduler]
command=php artisan schedule:work
directory=/var/www
user=www-data
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/logs/schedule.log
stderr_logfile=/var/www/storage/logs/schedule-error.log
priority=30
